# init.d/cicd-mcp
#!/sbin/openrc-run
# cicd-mcp: GitHub Actions CI/CD MCP server (FastAPI/Uvicorn :8012) の OpenRC 起動スクリプト
# 最小権限原則に従い llm-agent (非 root) ユーザで実行する

description="CI/CD (GitHub Actions) MCP Server (FastAPI/Uvicorn :8012)"

command="/opt/llm/venv/bin/uvicorn"
command_args="mcp.cicd.server:app --host 127.0.0.1 --port 8012 --workers 1"
command_background="yes"
# 専用非 root ユーザで起動する (sudo 権限なし)
command_user="llm-agent"

directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/cicd-mcp.log"
error_log="/opt/llm/logs/cicd-mcp.log"

# Load GITHUB_TOKEN from /etc/conf.d/cicd-mcp
rc_need="net"

depend() {
    need net
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
}
