# init.d/file-write-mcp
#!/sbin/openrc-run
# file-write-mcp: 書き込みファイル操作 MCP サーバ (FastAPI/Uvicorn :8007) の OpenRC 起動スクリプト

description="File Write MCP Server (FastAPI/Uvicorn :8007)"

command="/opt/llm/venv/bin/uvicorn"
command_args="file_write_mcp_server:app --host 127.0.0.1 --port 8007 --workers 1"
command_background="yes"
command_user="root"

directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/file-write-mcp.log"
error_log="/opt/llm/logs/file-write-mcp.log"

depend() {
    need net
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
}
