# init.d/file-delete-mcp
#!/sbin/openrc-run
# file-delete-mcp: 削除操作 MCP サーバ (FastAPI/Uvicorn :8008) の OpenRC 起動スクリプト

description="File Delete MCP Server (FastAPI/Uvicorn :8008)"

command="/opt/llm/venv/bin/uvicorn"
command_args="file_delete_mcp_server:app --host 127.0.0.1 --port 8008 --workers 1"
command_background="yes"
command_user="root"

directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/file-delete-mcp.log"
error_log="/opt/llm/logs/file-delete-mcp.log"

depend() {
    need net
}

start_pre() {
    # audit ログディレクトリを確保する
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
}
