# init.d/llama-agent
#!/sbin/openrc-run
# llama-agent: ローカル AI エージェント (CLI REPL) の OpenRC 起動スクリプト
# python -m agent で起動する対話型 REPL; uvicorn は使用しない

description="Local AI Agent (CLI REPL)"

command="/opt/llm/venv/bin/python"
command_args="-m agent"
command_background="yes"
command_user="root"

# scripts/ ディレクトリをカレントディレクトリとして起動する
# Python モジュール解決 (ConfigLoader 等の絶対インポート) はカレントディレクトリを参照するため必須
directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/agent.log"
error_log="/opt/llm/logs/agent.log"

depend() {
    # ネットワーク起動後に実行する
    need net
    # llama-server および web-search-mcp が起動してから実行する
    # file-mcp / github-mcp は agent が stdio サブプロセスとして起動するため依存不要
    after embed-llm llama-chat-llm llama-coding-llm web-search-mcp
}

start_pre() {
    # ログディレクトリが存在しない場合は作成する
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
}
