# init.d/llama-chat-llm
#!/sbin/openrc-run
# チャット LLM: gemma-4-e4b を port 8002 で起動する
# MQE (Multi-Query Expansion) と Cross-Encoder 再ランクにも兼用する

description="Chat LLM server (gemma-4-e4b) on :8002"

command="/opt/llm/llama.cpp/build/bin/llama-server"
command_args="
    --model /opt/llm/models/gemma-4-e4b-it-Q4_K_M.gguf
    --host 127.0.0.1
    --port 8002
    --ctx-size 4096
    --threads 4
    --n-predict 512
    --repeat-penalty 1.1"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/llama-chat-llm.log"
error_log="/opt/llm/logs/llama-chat-llm.log"

depend() {
    need net
}
