# init.d/rag-mcp
#!/sbin/openrc-run
# rag-mcp: 外部 RAG HTTP サービス (FastAPI + Uvicorn) の OpenRC 起動スクリプト

description="RAG MCP server (FastAPI/Uvicorn :8010)"

command="/opt/llm/venv/bin/uvicorn"
command_args="mcp.rag_pipeline.server:app --host 127.0.0.1 --port 8010 --workers 1"
command_background="yes"
command_user="root"

# rag_mcp_server.py が存在するディレクトリで起動する (Python モジュール解決のため)
directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/rag-mcp.log"
error_log="/opt/llm/logs/rag-mcp.log"

depend() {
    need net
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
}
