# init.d/shell-mcp
#!/sbin/openrc-run
# shell-mcp: サンドボックス付きシェル実行 MCP サーバ (FastAPI/Uvicorn :8009) の OpenRC 起動スクリプト
# 最小権限原則に従い llm-agent (非 root) ユーザで実行する

description="Shell Execution MCP Server (FastAPI/Uvicorn :8009)"

command="/opt/llm/venv/bin/uvicorn"
command_args="mcp.shell.server:app --host 127.0.0.1 --port 8009 --workers 1"
command_background="yes"
# 専用非 root ユーザで起動する (sudo 権限なし、/opt/llm/workspace 以外への書き込み権限なし)
command_user="llm-agent"

directory="/opt/llm/scripts"

pidfile="/run/${RC_SVCNAME}.pid"
output_log="/opt/llm/logs/shell-mcp.log"
error_log="/opt/llm/logs/shell-mcp.log"

depend() {
    need net
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /opt/llm/logs
    # workspace ディレクトリが存在しない場合は作成する (llm-agent が書き込み可能であること)
    checkpath --directory --owner llm-agent:llm-agent --mode 0755 /opt/llm/workspace
}
